# com.ai-news.healthcheck.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ai-news.healthcheck</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/chris/ai-news-dashboard/scraper/healthcheck.py</string>
        <string>--data</string>
        <string>/Users/chris/ai-news-dashboard/data.json</string>
        <string>--max-age</string>
        <string>30</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/chris/ai-news-dashboard</string>
    
    <!-- Check every 30 minutes -->
    <key>StartInterval</key>
    <integer>1800</integer>
    
    <key>StandardOutPath</key>
    <string>/tmp/ai-news-healthcheck.log</string>
    
    <key>StandardErrorPath</key>
    <string>/tmp/ai-news-healthcheck.err</string>
    
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
